# service/flask.conf
description "flask"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python2 /home/naboo/manual-aqbota/camera.py